# systemd/auto-debug.timer.template
[Unit]
Description=Claude Auto-Debug Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=%%INTERVAL%%
Persistent=true
Unit=auto-debug.service

[Install]
WantedBy=timers.target
